# deploy/deploy_consumer.service
[Unit]
Description=Cons1 Service
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
ExecStart=/bin/sh -c 'sleep 10; /usr/bin/php /home/mike/it490/cons1.php'
User=mike
Restart=on-failure


[Install]
WantedBy=multi-user.target
